# systemd.service
[Unit]
After=gamescope-session.service

[Service]
ExecStart=/usr/lib/hwsupport/wakehook
Restart=always
